# supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
directory=/usr/src/package/
command=gunicorn swagger_server.__main__:app -b 0.0.0.0:8000
autostart=true
autorestart=true
startsecs=10
startretries=10